# deploy/boya-agent.service
[Unit]
Description=BUAA Boya Course Push Agent
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/boya-agent
ExecStart=/home/ubuntu/boya-agent/venv/bin/python src/main.py
Restart=always
RestartSec=10
EnvironmentFile=/home/ubuntu/boya-agent/.env

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=boya-agent

[Install]
WantedBy=multi-user.target
